# deploy/campuslink-update.service
[Unit]
Description=Update CampusLink frontend from the GitHub deployment release
Wants=network-online.target
After=network-online.target nginx.service

[Service]
Type=oneshot
TimeoutStartSec=10min
ExecStart=/usr/local/sbin/campuslink-update
